# crontab
# Example of job definition:
0 * * * * /usr/local/bin/python /app/data_fetcher.py >> /var/log/cron.log 2>&1
#*/15 * * * * /usr/local/bin/python /app/data_fetcher.py >> /var/log/cron.log 2>&1
